# infra/systemd/vsa-agent.timer
[Unit]
Description=VSA Agent Timer — periodic sync every 30 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
